# magento/docker/containers/nginx/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:webserver]
command=/usr/sbin/nginx -g "daemon off;"